# .github/deb/frp/frpc.target
[Unit]
Description=Manage all the frpc instances
Requires=frpc-public.service frpc-local.service

[Install]
WantedBy=multi-user.target
