# .github/deb/frp/frps.target
[Unit]
Description=Manage all the frps instances
Requires=frps-local.service frps-public.service

[Install]
WantedBy=multi-user.target